# s6-rc.d/init-main/up
#!/command/execlineb -P
s6-setuidgid youruser
/opt/template-docker-python/init.sh
