# s6-rc.d/svc-main/run
#!/bin/sh
exec s6-setuidgid template-container-python uv run python3 -u /opt/template-container-python/main.py
