# s6-rc.d/svc-main/run
#!/usr/bin/env sh
uv run python3 -u /opt/template-docker-python/main.py
