# s6-rc.d/svc-main/run
#!/bin/sh
exec s6-setuidgid template-docker-python uv run python3 -u /opt/template-docker-python/main.py
