# s6-rc.d/svc-main/run
#!/usr/bin/env sh
# Service shell script goes here
